# configuration/supervisor_luigi.conf
[supervisord]
nodaemon = true

[program:luigi]
command = luigid
directory = /opt/code
stdout_logfile = /logs/luigi_out.log
stderr_logfile = /logs/luigi_err.log
process_name = %(program_name)s%(process_num)s
autorestart = true
