# configuration/supervisor_app.conf
[supervisord]
nodaemon = true

[program:datanectar]
command = /usr/local/bin/gunicorn app:app -w 2 -b 0.0.0.0:5000
directory = /opt/code
stdout_logfile = /logs/app_out.log
stderr_logfile = /logs/app_err.log
process_name = %(program_name)s%(process_num)s
autorestart = true
